# systemd-timer/myscript.timer
[Unit]
Description=mydescription

[Timer]
Unit=myscript.service
OnBootSec=2min

[Install]
WantedBy=timers.target